# blog.ini
[program:blogproject]
command= /usr/local/python3/bin/gunicorn blogproject.wsgi -c gunicorn_cnf.py
directory=/root/blogproject
autostart=true
autorestart=unexpected
user=root
stdout_logfile=/root/etc/supervisor/var/log/blog-tutorial-stdout.log
stderr_logfile=/root/etc/supervisor/var/log/blog-tutorial-stderr.log
